# br2-external-connect/board/stmicroelectronics/stm32mp157c-dk2/overlay/usr/lib/systemd/system/ofono.service
[Unit]
Description=Telephony service
After=syslog.target

[Service]
Type=dbus
BusName=org.ofono
ExecStartPre=-/etc/scripts/modem-reset
ExecStart=/usr/sbin/ofonod -n
StandardError=null

[Install]
WantedBy=multi-user.target
