# br2-external-adsb/board/orangepi/orangepi-zero/overlay/etc/systemd/system/adsb.service
[Unit]
Description=ADS-B
After=syslog.target network.target

[Service]
ExecStart=/bin/dump1090 --net --fix --aggressive --metric
KillMode=process
KillSignal=SIGKILL

